# ticketing-frontend.service
[Unit]
Description=Ticketing System Frontend
After=network.target ticketing-backend.service

[Service]
Type=simple
User=eazie
WorkingDirectory=/home/eazie/ticketing-system/frontend
ExecStartPre=/usr/bin/npm install
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npx serve -s build -l 3000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 